# deploy/question-bank.service
[Unit]
Description=题库系统 (gunicorn)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/srv/question-bank
EnvironmentFile=/etc/question-bank.env
ExecStart=/srv/question-bank/.venv/bin/gunicorn -w 2 -k gthread --threads 4 \
    --timeout 90 --graceful-timeout 30 \
    --error-logfile - \
    -b 127.0.0.1:8000 app:app
# 注:不开 --access-logfile——应用自带的 request JSON 日志(含 request_id/user_id/duration_ms)
# 是唯一访问日志源,避免双写与非 JSON 行污染日志流(Task 5 评审决定)
Restart=always
RestartSec=3

# 加固
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
